# stalwart-config/stalwart-apply-config.path
[Unit]

[Install]
WantedBy=multi-user.target

[Path]
PathExistsGlob=/run/stalwart_config/*.json
MakeDirectory=true
DirectoryMode=0700
