# stalwart-config/stalwart-apply-config.path
[Unit]
After=stalwart.service
Before=paths.target

[Path]
PathExistsGlob=/run/stalwart_config/*.json
MakeDirectory=/run/stalwart_config/
DirectoryMode=0700
