# stalwart-config/stalwart-apply-config.path
[Unit]
After=stalwart.service

[Install]
WantedBy=multi-user.target

[Path]
PathExistsGlob=/run/stalwart_config/*.json
MakeDirectory=/run/stalwart_config/
DirectoryMode=0700
