# set_hostname.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/set_hostname.sh

[Install]
WantedBy=multi-user.target

